# src/systemd/airdump.target
[Unit]
Description=Airdump Target
Documentation=https://github.com/yourrepo/airdump
After=network.target

[Install]
WantedBy=multi-user.target
